# systemd/owl-awdl.service
[Unit]
Description=OWL AWDL Interface (OpenDrop)
Documentation=https://owlink.org
After=network.target
Before=network-online.target

[Service]
Type=simple
RemainAfterExit=no
User=root
Group=root

# Create a virtual monitor interface on the WiFi adapter
ExecStartPre=/usr/sbin/iw dev wlo1 interface add mon0 type monitor
ExecStartPre=/usr/sbin/ip link set mon0 up
ExecStartPre=/usr/bin/sleep 1

# Start the OWL daemon on the monitor interface
# Note: OWL runs in foreground (Type=simple)
# If OWL crashes, check: owl -i mon0 output for errors
ExecStart=/usr/local/bin/owl -i mon0

# Clean up on stop: remove the virtual monitor interface
ExecStopPost=/usr/bin/sh -c 'iw dev mon0 del 2>/dev/null || true'
ExecStopPost=/usr/bin/sleep 1

# Bring WiFi interface back up (doesn't need password)
ExecStopPost=/usr/sbin/ip link set wlo1 up
ExecStopPost=/usr/bin/sleep 2

# Note: Manual WiFi reconnection required - see /usr/local/bin/wifi-reconnect.sh
# Automatic reconnection requires NetworkManager password which isn't available in systemd context

Restart=no
KillMode=process
StandardOutput=journal
StandardError=journal
SyslogIdentifier=owl-awdl
StandardInput=null

[Install]
WantedBy=multi-user.target
